# supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=/code/start_django.sh

[program:celery]
command=/code/start_supervisor.sh